# files/etc/profile.d/histtimeformat.sh
# Set the HISTTIMEFORMAT environment variable.

export HISTTIMEFORMAT="%d/%m/%y %T "
